# scripts/agriha-logger.service
[Unit]
Description=agriha history logger (MQTT -> SQLite)
After=network-online.target mosquitto.service
Wants=network-online.target
[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pi/agriha_logger.py
Restart=always
RestartSec=5
User=pi
[Install]
WantedBy=multi-user.target
